# Gateway/Raspbery/crontab
@reboot /home/pi/copyMotion/copy.sh >> /home/pi/copyMotion/logs/copy.log  2>&1